# systemd/otx-publisher.service
# /etc/systemd/system/otx-publisher.service
[Unit]
Description=Publish T-Pot IOCs to OTX
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/otx-publisher
ExecStart=/opt/otx-publisher/venv/bin/python /opt/otx-publisher/otx_tpot_publisher.py --config /opt/otx-publisher/config.json
StandardOutput=append:/opt/otx-publisher/run.log
StandardError=append:/opt/otx-publisher/run.log
