# systemd/otx-ciscoasa.service
[Unit]
Description=Cisco ASA → OTX hourly IP publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/otx-publisher
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /opt/otx-publisher/publisher/otx_ciscoasa_ips_only.py --config /opt/otx-publisher/config.ciscoasa.json
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
