# systemd/otx-publisher.timer
[Unit]
Description=Run OTX publisher every 6 hours
[Timer]
OnCalendar=*-*-* 00/6:17:00
Persistent=true
RandomizedDelaySec=120
[Install]
WantedBy=timers.target
